# upstart/mount_status_monitor.conf
description "Paranoid filesystem monitor"
author "Chris Adams <chris@improbable.org>"

start on runlevel [2345]
stop on runlevel [!2345]

console log

respawn

exec /usr/sbin/mount_status_monitor
